# docker/conf.d/chromium.conf
[program:chromium]
command=chromium --no-sandbox --remote-debugging-port=9222
autorestart=true
