# docker/conf.d/openbox.conf
[program:openbox]
command=bash /app/start-openbox.sh
environment=DISPLAY=:0
priority=10
autorestart=true
startsecs=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
